# run.service
[Uint]
After=network.service

[Service]
ExecStart=/usr/local/traing.sh

[Install]
WantedBy=multi-user.target
